# packages/spartan/extra/dcos-gen-resolvconf.timer
[Unit]
Description=Generate Resolv.conf Timer: Periodically update systemd-resolved for mesos-dns
[Timer]
OnBootSec=5sec
OnUnitActiveSec=1min
