# packages/spartan/extra/dcos-gen-resolvconf.service
[Unit]
Description=Generate Resolv.conf: Update systemd-resolved for MesosDNS
After=dcos-spartan.service

[Service]
Type=oneshot
StandardOutput=journal
StandardError=journal
EnvironmentFile=/etc/environment
EnvironmentFile=/opt/mesosphere/environment
EnvironmentFile=/opt/mesosphere/etc/dns_config
EnvironmentFile=/opt/mesosphere/etc/dns_search_config
EnvironmentFile=-/opt/mesosphere/etc/dns_config_master
ExecStartPre=/opt/mesosphere/active/toybox/bin/toybox timeout -k 1m 1m /bin/bash -c "! /opt/mesosphere/active/toybox/bin/toybox host ready.spartan"
ExecStartPre=/opt/mesosphere/active/toybox/bin/toybox timeout -k 1m 1m /opt/mesosphere/active/toybox/bin/toybox host ready.spartan 198.51.100.1
ExecStart=/opt/mesosphere/bin/gen_resolvconf.py /etc/resolv.conf
