# packages/spartan/extra/dcos-spartan-watchdog.timer
[Unit]
Description=Watchdog Spartan: Periodically check is Spartan is working
[Timer]
OnBootSec=5sec
OnUnitActiveSec=1min
