# systemd/dcos.target
[Unit]
After=dcos-bootstrap.service
Requires=dcos-bootstrap.service