# systemd/dcos.target
[Unit]
After=dcos-repair.service
Requires=dcos-repair.service
After=dcos-setup.service
Requires=dcos-setup.service