# packages/dcos-signal/extra/dcos-signal.timer
[Unit]
Description=Tracking every interval
[Timer]
OnCalendar=hourly
Unit=dcos-signal.service
