# packages/dcos-signal/extra/dcos-signal.timer
[Unit]
Description=Signal Timer: DCOS Signal Service Timer
[Timer]
OnCalendar=hourly
Unit=dcos-signal.service
