# packages/dcos-signal/extra/dcos-signal.service
[Unit]
Description=DCOS Tracking
After=dcos-mesos-master.service
[Service]
EnvironmentFile=/opt/mesosphere/environment
EnvironmentFile=-/opt/mesosphere/etc/cfn_signal_metadata
ExecStart=/opt/mesosphere/bin/dcos-signal --write_key=51ybGTeFEFU1xo6u10XMDrr6kATFyRyh
